# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/etc/service/apache2/run
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/apache2.log
stderr_logfile=/var/log/supervisor/apache2_err.log

[program:sogod]
command=/etc/service/sogod/run
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sogod.log
stderr_logfile=/var/log/supervisor/sogod_err.log

[program:memcached]
command=/etc/service/memcached/run
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/memcached.log
stderr_logfile=/var/log/supervisor/memcached_err.log